# ledserv.conf
[program:ledserv]
command=node /home/pi/ledserv/ledserv.js
autostart=true
autorestart=true
stderr_logfile=/var/log/ledserv.err.log
stdout_logfile=/var/log/ledserv.out.log
